# etc/skel/.config/autostart/a2jmidid.desktop
[Desktop Entry]
Exec=a2jmidid -e
Icon=system-run
Name=a2jmidid
Path=
Terminal=false
Type=Application
