# etc/skel/.config/autostart/cadence.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Name=Cadence Preload
Exec=cadence --minimized
Icon=cadence
Terminal=False
Type=Application
